# salt/restapi/restapi.cron

 30 * * * *  root  curl "http://control:8111/autoscale/?token=Haphohw6Efa9eiSu" >/dev/null 2>&1

